# private_Library/LaunchAgents/com.jon.topgrade.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.jon.topgrade</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/jon/.local/bin/run-topgrade-timer</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>3</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/jon/Library/Logs/topgrade/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jon/Library/Logs/topgrade/launchd.err.log</string>

    <key>RunAtLoad</key>
    <false/>
  </dict>
</plist>
